# systemd/proxmox-weekly-maintenance.service
[Unit]
Description=Proxmox VE Weekly Maintenance
Documentation=https://github.com/thelorax1775/proxmox-weekly-maintenance
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/proxmox-weekly-maintenance
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxmox-weekly-maintenance
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
